# deploy/site@.service
[Unit]
Description=customerenv site %i
After=ce-boot-config.service
[Service]
ExecStart=/usr/local/bin/http-nu /run/sites/%i.sock --log-format jsonl serve.nu
WorkingDirectory=/home/app/sites/%i/repo
EnvironmentFile=-/home/app/sites/%i/env
DynamicUser=yes
Group=sites
UMask=0007
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/app/sites/%i /run/sites
RuntimeDirectory=sites/%i-scratch
Restart=on-failure
[Install]
WantedBy=multi-user.target
